# index_devops.service
[Unit]
Description=Starting web services, index.js by nodejs
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /tmp/index.js
Restart=on-failure
